# systemd/user/mako-watcher.path
[Path]
PathModified=/home/xiangyi/.config/mako/config

[Install]
WantedBy=graphical-session.target
